# deploy/crea-dashboard-api.service
[Unit]
Description=CreaDashboard API server
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=__APP_DIR__/backend/.env
ExecStart=/usr/bin/env node dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crea-dashboard-api

[Install]
WantedBy=multi-user.target
